# roles/scalelite/files/prometheus-recordings.service
[Unit]
Description=Collect Scalelite Recordings Prometheus Metrics

[Service]
Type=oneshot
ExecStart=/usr/local/bin/scalelite-prometheus-recordings
