# roles/fluentd/files/td-agent-reload.service
[Unit]
Description=Restarts Td-Agent

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/tdagent-reload
