# roles/elasticsearch/files/espurge.service
[Unit]
Description=Purges ElasticSearch Older Indices

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/espurge
