# roles/bigbluebutton/files/bbb-services.timer
[Unit]
Description=Run BigBlueButton Services Prometheus Collector

[Timer]
OnBootSec=5min
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
