# roles/bigbluebutton/files/recordings-purge.timer
[Unit]
Description=Run BigBlueButton Recordings Purge

[Timer]
OnBootSec=90min
OnCalendar=0/6:00:00

[Install]
WantedBy=timers.target
