# roles/bigbluebutton/files/recordings-purge.service
[Unit]
Description=Purges Older BigBlueButton Recordings

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/bigbluebutton-recordings-purge
