# roles/commons/files/prometheus-check_certificates.service
[Unit]
Description=Collect x509 Certificates TTL Metrics

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/certificates-prometheus
