# roles/commons/files/prometheus-check_certificates.timer
[Unit]
Description=Run x509 Certificates TTL Collector

[Timer]
OnBootSec=55min
OnCalendar=0/3:00:00

[Install]
WantedBy=timers.target
