# roles/fail2ban/files/prometheus-check_fail2ban.service
[Unit]
Description=Collect Fail2ban Metrics

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/fail2ban-prometheus
