# roles/fail2ban/files/prometheus-check_fail2ban.timer
[Unit]
Description=Run Fail2ban Metrics Collector

[Timer]
OnBootSec=5min
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
